# deploy/secure-mcp.service
# systemd unit for secure-mcp.
#
# Install as /etc/systemd/system/secure-mcp@.service so multiple instances can
# coexist (one per identity file): `systemctl start secure-mcp@soc-analyst`.
#
# Edit ExecStart to point at your wrapper script that loads secrets from
# Vault (or equivalent) into env before exec-ing the entrypoint. Static env
# files in /etc are NOT recommended — secrets should never sit on disk.

[Unit]
Description=Secure MCP broker (Check Point TE + Lakera Guard) — instance %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=secure-mcp
Group=secure-mcp
WorkingDirectory=/opt/secure-mcp

# Wrapper is responsible for: (1) pulling secrets from Vault into env,
# (2) exec-ing the real binary. The wrapper takes the identity-file path as
# its argument so each instance binds to a distinct caller.
ExecStart=/opt/secure-mcp/deploy/vault-env-load.sh /etc/secure-mcp/identity-%i.json /opt/secure-mcp/.venv/bin/secure-mcp

Restart=on-failure
RestartSec=5s

# --- Hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=

# Allow writes only to the audit log, the staging dir, and runtime cache.
ReadWritePaths=/var/log/secure-mcp /var/lib/secure-mcp

# Limit outbound network reachability to public IPs only (defense in depth
# against accidental loopback / private-net targeting).
IPAddressDeny=any
IPAddressAllow=0.0.0.0/0
IPAddressDeny=10.0.0.0/8 172.16.0.0/12 192.168.0.0/16 127.0.0.0/8 169.254.0.0/16

[Install]
WantedBy=multi-user.target
